# scripts/vm/idle-shutdown.service
[Unit]
Description=ComfyUI Idle Shutdown
After=network.target comfyui.service
Wants=comfyui.service

[Service]
Type=simple
User=renato
Environment=IDLE_MINUTES=30
ExecStart=/opt/ai-studio/idle-shutdown.sh
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
